# ansible.timer
[Unit]
Description=Ansible

[Timer]
OnStartupSec=5min
OnUnitActiveSec=30min
RandomizedDelaySec=5min

[Install]
WantedBy=timers.target
